# dist/supervisor/wserver.conf
[program:wserver]
command=/usr/bin/wserver
autostart=true
autorestart=true
